# install/etc/cron.daily/couch-db
#!/bin/bash

bash /root/thinx-device-api/couch_cleanup.sh 
