# install/etc/cron.d/repoclean
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
0 */12 * * * root cd /mnt/local/repos && bash ./cleanup.sh